# systemd/certheim-expiry-warn.timer
[Unit]
Description=Daily Certheim expiry warning run

[Timer]
OnCalendar=*-*-* 06:30:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
