# systemd/certinel-doctor.service
[Unit]
# Periodic Certinel health check (certinel-doctor) that emails on failure via
# Mailgun. Config: /etc/certinel/doctor-alert.conf. Runs as root because the
# doctor probes service-account traversal + SELinux labels. Driven by the timer.
Description=Certinel health check + alert on failure
After=certinel-api.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/certinel-doctor-alert
