# systemd/certinel-tls-renew.timer
[Unit]
Description=Daily renewal of the Certinel nginx TLS certificate

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
