# systemd/certinel-deliver.service
[Unit]
Description=Certinel - deliver issued certificates to their destinations (retry pass)
After=network-online.target

[Service]
Type=oneshot
User=certinel
Group=certinel
WorkingDirectory=/opt/certinel
EnvironmentFile=-/etc/certinel/certinel.env
ExecStart=/opt/certinel/venv/bin/python3 -c 'import app; r = app.run_deliveries(); print(f"delivery: {r}")'
Environment=PYTHONUNBUFFERED=1
