# systemd/certinel-expiry-warn.service
[Unit]
Description=Certinel - send certificate expiry warnings
After=network-online.target

[Service]
Type=oneshot
User=csrapi
Group=csrapi
WorkingDirectory=/opt/csr-dashboard
EnvironmentFile=-/etc/csr-dashboard/csr-dashboard.env
ExecStart=/opt/csr-dashboard/venv/bin/python3 -c 'import app; s,e = app.run_expiry_warnings(); print(f"expiry warnings: sent={s} errors={e}")'
Environment=PYTHONUNBUFFERED=1
