# systemd/certinel-auto-renew.service
[Unit]
Description=Certinel - automatically renew certificates nearing expiry
After=network-online.target

[Service]
Type=oneshot
User=certinel
Group=certinel
WorkingDirectory=/opt/certinel
EnvironmentFile=-/etc/certinel/certinel.env
ExecStart=/opt/certinel/venv/bin/python3 -c 'import app; r,s,e = app.run_auto_renew(); print(f"auto-renew: renewed={r} skipped={s} errors={e}")'
Environment=PYTHONUNBUFFERED=1
