# systemd/certinel-doctor.timer
[Unit]
Description=Run the Certheim health check every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
