# systemd/certinel-slack-listener.service
[Unit]
# Homelab-only Slack Socket Mode listener. Outbound-only (no inbound exposure).
# NOT installed by deploy.sh - deploy + enable manually where wanted. Air-gapped
# prod has no path to Slack, so leave it off there.
Description=Certinel Slack Socket Mode listener
After=network-online.target certinel-api.service
Wants=network-online.target

[Service]
User=certinel
Group=certinel
EnvironmentFile=-/etc/certinel/certinel.env
ExecStart=/opt/certinel/venv/bin/python /opt/certinel/slack_listener.py
Restart=always
RestartSec=10
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
