# systemd/certheim-slack-listener.service
[Unit]
# Homelab-only Slack Socket Mode listener. Outbound-only (no inbound exposure).
# NOT installed by deploy.sh - deploy + enable manually where wanted. Air-gapped
# prod has no path to Slack, so leave it off there.
Description=Certheim Slack Socket Mode listener
After=network-online.target certheim-api.service
Wants=network-online.target

[Service]
User=certheim
Group=certheim
EnvironmentFile=-/etc/certheim/certheim.env
ExecStart=/opt/certheim/venv/bin/python /opt/certheim/slack_listener.py
Restart=always
RestartSec=10
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
