# systemd/csr-deliver.service
[Unit]
Description=CSR Dashboard - deliver issued certificates to their destinations (retry pass)
After=network-online.target

[Service]
Type=oneshot
User=csrapi
Group=csrapi
WorkingDirectory=/opt/csr-dashboard
EnvironmentFile=-/etc/csr-dashboard/csr-dashboard.env
ExecStart=/opt/csr-dashboard/venv/bin/python3 -c 'import app; r = app.run_deliveries(); print(f"delivery: {r}")'
Environment=PYTHONUNBUFFERED=1
